# debian/prerm
#!/bin/sh

set -e

if [ -x "/etc/init.d/dtc-xen" ] ; then
	if [ -x /usr/sbin/invoke-rc.d ] ; then
		/usr/sbin/invoke-rc.d dtc-xen stop
	else
		/etc/init.d/dtc-xen stop
	fi
fi

#DEBHELPER#

exit 0
